# base-image/etc/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info
logfile_maxbytes=16MB
logfile_backups=2
user=root
logfile=/srv/asiou/log/supervisord.log
childlogdir=/srv/asiou/log/

[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
autostart=true
autorestart=true
user=nobody
priority=90
stdout_logfile=/srv/asiou/log/redis-stdout.log
stderr_logfile=/srv/asiou/log/redis-stderr.log

[program:apache2-pre]
command=/bin/sh -c 'rm /var/run/apache2/*; :'
autorestart=false
exitcodes=0
startsecs=0
startretries=0
priority=99

[program:apache2]
command=/usr/sbin/apachectl -D FOREGROUND
autostart=true
autorestart=true
user=www-data
priority=100
stdout_logfile=/srv/asiou/log/apache2-stdout.log
stderr_logfile=/srv/asiou/log/apache2-stderr.log

[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true
user=www-data
priority=200
stdout_logfile=/srv/asiou/log/nginx-stdout.log
stderr_logfile=/srv/asiou/log/nginx-stderr.log

[program:cron]
command=/usr/sbin/cron -f -L 15
autostart=true
autorestart=true
user=root
priority=300
stdout_logfile=/srv/asiou/log/cron-stdout.log
stderr_logfile=/srv/asiou/log/cron-stderr.log
